# ch11-deployment/listing-11_1/hellonode.conf
author      "Robert DeGrimston"
description "hellonode"
setuid      "nonrootuser"
start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown
respawn
console log
env NODE_ENV=production
exec /usr/bin/node /path/to/server.js
